# esp-sync.path
[Path]
# Trigger once per boot when ostree finalises the new deployment
Unit=esp-sync.service
PathChanged=/run/ostree-booted