# sd-bootc.path
[Path]
# Trigger once per boot when ostree finalises the new deployment
Unit=sd-bootc.service
PathChanged=/run/ostree-booted